# docs/goop-rendezvous.service
[Unit]
Description=Goop² Rendezvous Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=goop
Group=goop
WorkingDirectory=/opt/goop

# Path to the goop2 binary in rendezvous mode
ExecStart=/opt/goop/goop2 -rendezvous -addr 127.0.0.1:8787

# Restart policy
Restart=always
RestartSec=10

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/goop

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=goop-rendezvous

[Install]
WantedBy=multi-user.target
